# data/tallow.service.in

[Unit]
Description=Tallow Service

[Service]
ExecStart=@prefix@/sbin/tallow
Restart=always
RestartSec=3
Nice=10

[Install]
WantedBy=network.target

